# src/dotfiles/config/systemd/user/multi-user.target.wants/torrents-for-transmission.path
[Unit]
Description=Watch for downloaded torrent files

[Path]
PathChanged=%h/Downloads/torrent
Unit=torrents-for-transmission.service

[Install]
WantedBy=multi-user.target
